# branch/testing/etc/cron.daily/clean-alarms
#!/bin/sh
/usr/bin/php /usr/share/ossim/scripts/alarms/delete_alarms.php
